# supervisord.conf
[supervisord]
nodaemon=true

[program:mongodb]
command=/usr/bin/mongod --smallfiles
autorestart=true

[program:walle]
command=/usr/bin/nodejs /opt/wall-e/server.js
autorestart=true
